# config/supervisor-squash-tm.ini
[program:squash-tm]
directory=%(ENV_SQUASH_TM_DIR)s/bin
command=./startup.sh
priority=1
autostart=true
autorestart=true
stdout_logfile=%(ENV_SQUASH_TM_STDOUT_FILE)s
stderr_logfile=%(ENV_SQUASH_TM_STDERR_FILE)s

